# playbooks/templates/jellyfin.container.j2
[Unit]
Description=Jellyfin Media Server
Description=Jellyfin Media Server
After=network-online.target
Wants=network-online.target

# Dépendances plus souples (seulement After, pas Requires)
After=var-home-famille-nas-video.mount var-home-famille-nas-PlexMediaServer.mount

[Container]
Image=docker.io/jellyfin/jellyfin:latest
ContainerName=jellyfin

# Config et cache
Volume={{ jellyfin_config_dir }}:/config:Z
Volume={{ jellyfin_cache_dir }}:/cache:Z

# Médiathèques NAS
Volume=/var/home/{{ ansible_user }}/nas/video:/media/video:ro
Volume=/var/home/{{ ansible_user }}/nas/PlexMediaServer:/media/plex:ro

# Disque interne
Volume=/var/home/{{ ansible_user }}/internes/exos/Videos:/media/exos:ro

Network=host
PublishPort=127.0.0.1:8096:8096

# Performance
Security=seccomp=unconfined
Environment=TZ=Europe/Paris

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=default.target