# playbooks/templates/jellyfin.container.j2
[Unit]
Description=Jellyfin Media Server
After=var-home-famille-nas-PlexMediaServer.mount var-home-famille-nas-video.mount
Requires=var-home-famille-nas-PlexMediaServer.mount var-home-famille-nas-video.mount
[Container]
Image=docker.io/jellyfin/jellyfin:latest
ContainerName=jellyfin

# Données et config Jellyfin (persistantes)
Volume=/var/home/{{ ansible_user }}/.local/share/jellyfin/config:/config:Z
Volume=/var/home/{{ ansible_user }}/.local/share/jellyfin/cache:/cache:Z

# Montages NAS en lecture seule
Volume=/var/home/{{ ansible_user }}/nas/video:/media/video
Volume=/var/home/{{ ansible_user }}/nas/PlexMediaServer:/media/plex

# Disque interne
Volume=/var/home/{{ ansible_user }}/internes/exos/Videos:/media/exos

# Réseau
Network=host
PublishPort=127.0.0.1:8096:8096

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target